# examples/cu_rp_balancebot/etc/99-pwm-permissions.rules
SUBSYSTEM=="pwm", RUN+="/usr/local/bin/set_pwm_permissions.sh"
